# ssclient/crontab
# ssclient crontab
*/1 * * * *    shankar  /home/shankar/src/ssclient/ssclient.py
